# omero-in-a-box-wndcharm/omero-wndcharm.conf
[program:omero-wndcharm]
command=/omero-wndcharm
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=510
